# garmy_mcp/rootfs/etc/services.d/garmy_sync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Periodic Garmy sync service
# ==============================================================================

export HOME="/data"

db_path="$(bashio::config 'database.path')"
interval_hours="$(bashio::config 'sync.interval_hours')"
sync_days="$(bashio::config 'sync.history_days')"

export GARMY_DB_PATH="${db_path}"
export GARMY_SYNC_DAYS="${sync_days}"

if ! [[ "${interval_hours}" =~ ^[0-9]+$ ]]; then
    bashio::log.warning "sync.interval_hours ('${interval_hours}') is not a valid whole number; defaulting to 6."
    interval_hours=6
fi
interval_seconds=$(( interval_hours * 3600 ))

while true; do
    bashio::log.info "Starting Garmy sync..."
    garmy-sync-ha
    code=$?
    if [[ ${code} -eq 0 ]]; then
        bashio::log.info "Garmy sync completed successfully."
    elif [[ ${code} -eq 2 ]]; then
        bashio::log.warning "Garmy sync skipped: no tokens available."
        bashio::log.warning "Set garmin.oauth1_token_json/oauth2_token_json in config, or copy token files to /data/.garmy/."
    else
        bashio::log.warning "Garmy sync finished with errors (exit ${code})."
    fi
    bashio::log.info "Next sync in ${interval_hours} hour(s)."
    sleep "${interval_seconds}" &
    wait $!
done
