# google_workspace_mcp/rootfs/etc/services.d/google_workspace_mcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Google Workspace MCP service
# ==============================================================================

set -e

declare -a cmd
declare -a tools
declare -a permissions

normalize_optional() {
    local value="$1"
    local lowered
    lowered="$(echo "${value}" | tr '[:upper:]' '[:lower:]')"

    case "${lowered}" in
        ""|"null"|"none")
            echo ""
            ;;
        *)
            echo "${value}"
            ;;
    esac
}

auth_client_id="$(bashio::config 'google_oauth_client_id')"
auth_client_secret="$(bashio::config 'google_oauth_client_secret')"
credentials_file="$(bashio::config 'authentication.google_client_secret_path')"
[[ -z "${credentials_file}" ]] && credentials_file="$(bashio::config 'google_client_secret_path')"
user_google_email="$(bashio::config 'authentication.user_google_email')"
[[ -z "${user_google_email}" ]] && user_google_email="$(bashio::config 'user_google_email')"
tool_tier="$(bashio::config 'tool_selection.workspace_mcp_tool_tier')"
[[ -z "${tool_tier}" ]] && tool_tier="$(bashio::config 'tool_tier')"

credentials_dir="$(bashio::config 'credential_store.workspace_mcp_credentials_dir')"
[[ -z "${credentials_dir}" ]] && credentials_dir="$(bashio::config 'workspace_mcp_credentials_dir')"
[[ -z "${credentials_dir}" ]] && credentials_dir="$(bashio::config 'advanced.workspace_mcp_credentials_dir')"

base_uri="$(bashio::config 'server.workspace_mcp_base_uri')"
[[ -z "${base_uri}" ]] && base_uri="$(bashio::config 'advanced.server.workspace_mcp_base_uri')"
[[ -z "${base_uri}" ]] && base_uri="$(bashio::config 'workspace_mcp_base_uri')"

# Prefer flattened top-level groups and fall back to older advanced/flat keys.
[[ -z "${credentials_file}" ]] && credentials_file="$(bashio::config 'advanced.oauth.google_client_secret_path')"
[[ -z "${credentials_file}" ]] && credentials_file="$(bashio::config 'advanced.google_client_secret_path')"
[[ -z "${credentials_file}" ]] && credentials_file="$(bashio::config 'google_client_secret_path')"
external_url="$(bashio::config 'workspace_external_url')"
[[ -z "${external_url}" ]] && external_url="$(bashio::config 'mcp_host.workspace_external_url')"
[[ -z "${external_url}" ]] && external_url="$(bashio::config 'oauth.workspace_external_url')"
[[ -z "${external_url}" ]] && external_url="$(bashio::config 'advanced.oauth.workspace_external_url')"
[[ -z "${external_url}" ]] && external_url="$(bashio::config 'advanced.workspace_external_url')"
redirect_uri="$(bashio::config 'oauth21.google_oauth_redirect_uri')"
[[ -z "${redirect_uri}" ]] && redirect_uri="$(bashio::config 'oauth.google_oauth_redirect_uri')"
[[ -z "${redirect_uri}" ]] && redirect_uri="$(bashio::config 'advanced.oauth.google_oauth_redirect_uri')"
[[ -z "${redirect_uri}" ]] && redirect_uri="$(bashio::config 'advanced.google_oauth_redirect_uri')"
[[ -z "${redirect_uri}" ]] && redirect_uri="$(bashio::config 'google_oauth_redirect_uri')"
custom_redirect_uris="$(bashio::config 'oauth21.oauth_custom_redirect_uris')"
[[ -z "${custom_redirect_uris}" ]] && custom_redirect_uris="$(bashio::config 'oauth.oauth_custom_redirect_uris')"
[[ -z "${custom_redirect_uris}" ]] && custom_redirect_uris="$(bashio::config 'advanced.oauth.oauth_custom_redirect_uris')"
[[ -z "${custom_redirect_uris}" ]] && custom_redirect_uris="$(bashio::config 'advanced.oauth_custom_redirect_uris')"
[[ -z "${custom_redirect_uris}" ]] && custom_redirect_uris="$(bashio::config 'oauth_custom_redirect_uris')"
allowed_origins="$(bashio::config 'oauth21.oauth_allowed_origins')"
[[ -z "${allowed_origins}" ]] && allowed_origins="$(bashio::config 'oauth.oauth_allowed_origins')"
[[ -z "${allowed_origins}" ]] && allowed_origins="$(bashio::config 'advanced.oauth.oauth_allowed_origins')"
[[ -z "${allowed_origins}" ]] && allowed_origins="$(bashio::config 'advanced.oauth_allowed_origins')"
[[ -z "${allowed_origins}" ]] && allowed_origins="$(bashio::config 'oauth_allowed_origins')"
jwt_signing_key="$(bashio::config 'oauth21.fastmcp_server_auth_google_jwt_signing_key')"
[[ -z "${jwt_signing_key}" ]] && jwt_signing_key="$(bashio::config 'oauth.fastmcp_server_auth_google_jwt_signing_key')"
[[ -z "${jwt_signing_key}" ]] && jwt_signing_key="$(bashio::config 'advanced.oauth.fastmcp_server_auth_google_jwt_signing_key')"
[[ -z "${jwt_signing_key}" ]] && jwt_signing_key="$(bashio::config 'advanced.fastmcp_server_auth_google_jwt_signing_key')"
[[ -z "${jwt_signing_key}" ]] && jwt_signing_key="$(bashio::config 'fastmcp_server_auth_google_jwt_signing_key')"
transport="$(bashio::config 'mcp_host.transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'server.workspace_mcp_transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'advanced.server.workspace_mcp_transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'advanced.mcp_host.transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'advanced.server.transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'advanced.transport')"
[[ -z "${transport}" ]] && transport="$(bashio::config 'transport')"
selected_tools="$(bashio::config 'tool_selection.workspace_mcp_tools')"
[[ -z "${selected_tools}" ]] && selected_tools="$(bashio::config 'access.tools')"
[[ -z "${selected_tools}" ]] && selected_tools="$(bashio::config 'advanced.access.tools')"
[[ -z "${selected_tools}" ]] && selected_tools="$(bashio::config 'advanced.tools')"
[[ -z "${selected_tools}" ]] && selected_tools="$(bashio::config 'tools')"
permissions_raw="$(bashio::config 'tool_selection.workspace_mcp_permissions')"
[[ -z "${permissions_raw}" ]] && permissions_raw="$(bashio::config 'access.permissions')"
[[ -z "${permissions_raw}" ]] && permissions_raw="$(bashio::config 'advanced.access.permissions')"
[[ -z "${permissions_raw}" ]] && permissions_raw="$(bashio::config 'advanced.permissions')"
[[ -z "${permissions_raw}" ]] && permissions_raw="$(bashio::config 'permissions')"
allowed_file_dirs="$(bashio::config 'local_files.allowed_file_dirs')"
[[ -z "${allowed_file_dirs}" ]] && allowed_file_dirs="$(bashio::config 'access.allowed_file_dirs')"
[[ -z "${allowed_file_dirs}" ]] && allowed_file_dirs="$(bashio::config 'advanced.access.allowed_file_dirs')"
[[ -z "${allowed_file_dirs}" ]] && allowed_file_dirs="$(bashio::config 'advanced.allowed_file_dirs')"
[[ -z "${allowed_file_dirs}" ]] && allowed_file_dirs="$(bashio::config 'allowed_file_dirs')"
port="$(bashio::config 'mcp_host.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'server.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'advanced.server.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'advanced.mcp_host.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'advanced.server.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'advanced.workspace_mcp_port')"
[[ -z "${port}" ]] && port="$(bashio::config 'workspace_mcp_port')"
host="$(bashio::config 'mcp_host.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'server.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'advanced.server.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'advanced.mcp_host.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'advanced.server.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'advanced.workspace_mcp_host')"
[[ -z "${host}" ]] && host="$(bashio::config 'workspace_mcp_host')"

oauth_storage_backend="$(bashio::config 'oauth21.workspace_mcp_oauth_proxy_storage_backend')"
[[ -z "${oauth_storage_backend}" ]] && oauth_storage_backend="$(bashio::config 'advanced.oauth.workspace_mcp_oauth_proxy_storage_backend')"
[[ -z "${oauth_storage_backend}" ]] && oauth_storage_backend="$(bashio::config 'workspace_mcp_oauth_proxy_storage_backend')"

allowed_client_redirect_uris="$(bashio::config 'oauth21.workspace_mcp_allowed_client_redirect_uris')"
[[ -z "${allowed_client_redirect_uris}" ]] && allowed_client_redirect_uris="$(bashio::config 'advanced.oauth.workspace_mcp_allowed_client_redirect_uris')"
[[ -z "${allowed_client_redirect_uris}" ]] && allowed_client_redirect_uris="$(bashio::config 'workspace_mcp_allowed_client_redirect_uris')"

auth_client_id="$(normalize_optional "${auth_client_id}")"
auth_client_secret="$(normalize_optional "${auth_client_secret}")"
credentials_file="$(normalize_optional "${credentials_file}")"
user_google_email="$(normalize_optional "${user_google_email}")"
tool_tier="$(normalize_optional "${tool_tier}")"
credentials_dir="$(normalize_optional "${credentials_dir}")"
base_uri="$(normalize_optional "${base_uri}")"
external_url="$(normalize_optional "${external_url}")"
redirect_uri="$(normalize_optional "${redirect_uri}")"
custom_redirect_uris="$(normalize_optional "${custom_redirect_uris}")"
allowed_origins="$(normalize_optional "${allowed_origins}")"
jwt_signing_key="$(normalize_optional "${jwt_signing_key}")"
transport="$(normalize_optional "${transport}")"
selected_tools="$(normalize_optional "${selected_tools}")"
permissions_raw="$(normalize_optional "${permissions_raw}")"
allowed_file_dirs="$(normalize_optional "${allowed_file_dirs}")"
port="$(normalize_optional "${port}")"
host="$(normalize_optional "${host}")"
oauth_storage_backend="$(normalize_optional "${oauth_storage_backend}")"
allowed_client_redirect_uris="$(normalize_optional "${allowed_client_redirect_uris}")"

[[ -z "${transport}" ]] && transport="streamable-http"
[[ -z "${host}" ]] && host="0.0.0.0"
[[ -z "${port}" ]] && port="8000"
[[ -z "${allowed_file_dirs}" ]] && allowed_file_dirs="/share"
[[ -z "${credentials_dir}" ]] && credentials_dir="/data/google_workspace_credentials"
[[ -z "${oauth_storage_backend}" ]] && oauth_storage_backend="disk"

mkdir -p "${credentials_dir}"
export GOOGLE_MCP_CREDENTIALS_DIR="${credentials_dir}"
export WORKSPACE_MCP_CREDENTIALS_DIR="${credentials_dir}"
export WORKSPACE_MCP_PORT="${port}"
export WORKSPACE_MCP_HOST="${host}"
export ALLOWED_FILE_DIRS="${allowed_file_dirs}"
export WORKSPACE_MCP_OAUTH_PROXY_STORAGE_BACKEND="${oauth_storage_backend}"

[[ -n "${auth_client_id}" ]] && export GOOGLE_OAUTH_CLIENT_ID="${auth_client_id}"
[[ -n "${auth_client_secret}" ]] && export GOOGLE_OAUTH_CLIENT_SECRET="${auth_client_secret}"
[[ -n "${credentials_file}" ]] && export GOOGLE_CLIENT_SECRET_PATH="${credentials_file}"
[[ -n "${user_google_email}" ]] && export USER_GOOGLE_EMAIL="${user_google_email}"
[[ -n "${external_url}" ]] && export WORKSPACE_EXTERNAL_URL="${external_url}"
[[ -n "${base_uri}" ]] && export WORKSPACE_MCP_BASE_URI="${base_uri}"
[[ -n "${redirect_uri}" ]] && export GOOGLE_OAUTH_REDIRECT_URI="${redirect_uri}"
[[ -n "${custom_redirect_uris}" ]] && export OAUTH_CUSTOM_REDIRECT_URIS="${custom_redirect_uris}"
[[ -n "${allowed_origins}" ]] && export OAUTH_ALLOWED_ORIGINS="${allowed_origins}"
[[ -n "${jwt_signing_key}" ]] && export FASTMCP_SERVER_AUTH_GOOGLE_JWT_SIGNING_KEY="${jwt_signing_key}"
[[ -n "${allowed_client_redirect_uris}" ]] && export WORKSPACE_MCP_ALLOWED_CLIENT_REDIRECT_URIS="${allowed_client_redirect_uris}"

if bashio::config.true 'authentication.oauthlib_insecure_transport' || bashio::config.true 'oauth.oauthlib_insecure_transport' || bashio::config.true 'advanced.oauth.oauthlib_insecure_transport' || bashio::config.true 'advanced.oauthlib_insecure_transport' || bashio::config.true 'oauthlib_insecure_transport'; then
    export OAUTHLIB_INSECURE_TRANSPORT=1
fi
if bashio::config.true 'oauth.mcp_enable_oauth21' || bashio::config.true 'advanced.oauth.mcp_enable_oauth21' || bashio::config.true 'advanced.mcp_enable_oauth21' || bashio::config.true 'mcp_enable_oauth21'; then
    export MCP_ENABLE_OAUTH21=true
fi
if bashio::config.true 'oauth.external_oauth21_provider' || bashio::config.true 'advanced.oauth.external_oauth21_provider' || bashio::config.true 'advanced.external_oauth21_provider' || bashio::config.true 'external_oauth21_provider'; then
    export EXTERNAL_OAUTH21_PROVIDER=true
fi
if bashio::config.true 'mcp_host.workspace_mcp_stateless_mode' || bashio::config.true 'advanced.mcp_host.workspace_mcp_stateless_mode' || bashio::config.true 'advanced.server.workspace_mcp_stateless_mode' || bashio::config.true 'advanced.workspace_mcp_stateless_mode' || bashio::config.true 'workspace_mcp_stateless_mode'; then
    export WORKSPACE_MCP_STATELESS_MODE=true
fi
if bashio::config.true 'oauth21.mcp_enable_oauth21'; then
    export MCP_ENABLE_OAUTH21=true
fi
if bashio::config.true 'oauth21.external_oauth21_provider'; then
    export EXTERNAL_OAUTH21_PROVIDER=true
fi
if bashio::config.true 'oauth21.workspace_mcp_stateless_mode'; then
    export WORKSPACE_MCP_STATELESS_MODE=true
fi

cmd=(workspace-mcp --transport "${transport}")

if [[ -n "${tool_tier}" ]]; then
    cmd+=(--tool-tier "${tool_tier}")
fi

if [[ -n "${selected_tools}" ]]; then
    selected_tools="${selected_tools//,/ }"
    read -r -a tools <<< "${selected_tools}"
    if [[ ${#tools[@]} -gt 0 ]]; then
        cmd+=(--tools "${tools[@]}")
    fi
fi

if [[ -n "${permissions_raw}" ]]; then
    permissions_raw="${permissions_raw//,/ }"
    read -r -a permissions <<< "${permissions_raw}"
    if [[ ${#permissions[@]} -gt 0 ]]; then
        cmd+=(--permissions "${permissions[@]}")
    fi
fi

if bashio::config.true 'tool_selection.workspace_mcp_read_only' || bashio::config.true 'access.read_only' || bashio::config.true 'advanced.access.read_only' || bashio::config.true 'advanced.read_only' || bashio::config.true 'read_only'; then
    cmd+=(--read-only)
fi
if bashio::config.true 'authentication.single_user' || bashio::config.true 'access.single_user' || bashio::config.true 'advanced.access.single_user' || bashio::config.true 'advanced.single_user' || bashio::config.true 'single_user'; then
    cmd+=(--single-user)
fi

bashio::log.info "Starting Google Workspace MCP on ${host}:${port} with transport '${transport}'."

exec "${cmd[@]}"
