# garmy_mcp/rootfs/etc/services.d/garmy_mcp/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the Garmy MCP service
# ==============================================================================

set -euo pipefail

if bashio::config.true 'advanced.debug_logging'; then
    bashio::log.level "debug"
    export FASTMCP_LOG_LEVEL="DEBUG"
    set -x
fi

db_path="$(bashio::config 'database.path')"
# Fixed: FastMCP's HTTP transport is the only one this add-on supports.
transport="streamable-http"
host="$(bashio::config 'mcp.host')"
# Must match the container-side port declared in config.yaml's `ports:` map;
# it is not user-configurable because Supervisor's port publishing is static.
port="8000"
path="$(bashio::config 'mcp.path')"

declare -a cmd

if [[ -z "${db_path}" ]]; then
    bashio::log.fatal "database.path is required"
    bashio::exit.nok
fi

mkdir -p "$(dirname "${db_path}")"

if [[ ! -f "${db_path}" ]]; then
    bashio::log.warning "Database file not found at ${db_path}; creating an empty SQLite database."
    python3 -c "import sqlite3, sys; sqlite3.connect(sys.argv[1]).close()" "${db_path}"
fi

if [[ -z "${host}" ]]; then
    bashio::log.fatal "mcp.host is required."
    bashio::exit.nok
fi

if [[ -z "${path}" ]]; then
    bashio::log.fatal "mcp.path is required."
    bashio::exit.nok
fi

export FASTMCP_TRANSPORT="${transport}"
export FASTMCP_HOST="${host}"
export FASTMCP_PORT="${port}"
export FASTMCP_STREAMABLE_HTTP_PATH="${path}"

cmd=(
    garmy-mcp
    server
    --database "${db_path}"
)

exec "${cmd[@]}"
